# Ansible/roles/etcd_cluster/files/etcd.service
[Unit]
Description=ETCD cluster for Kubernetes

[Service]
Type=notify
ExecStart=/etc/etcd/startup
Restart=always